# veriqko-api.service
[Unit]
Description=Veriqko API Server
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=veriqko
Group=veriqko
WorkingDirectory=/opt/veriqko/app/apps/api
Environment=PATH=/opt/veriqko/app/apps/api/venv/bin:/usr/bin
ExecStart=/opt/veriqko/app/apps/api/venv/bin/uvicorn veriqko.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=append:/opt/veriqko/logs/api.log
StandardError=append:/opt/veriqko/logs/api-error.log

[Install]
WantedBy=multi-user.target
